# vdr2.service.example
.include /lib/systemd/system/vdr.service
# See /usr/share/doc/vdr/README.*.multi-vdr.

[Unit]
Description=Video Disk Recorder instance 2

[Service]
Environment=VDR_INSTANCE=vdr2
